# debian/jlcarveth-dev/DEBIAN/postinst
#!/bin/bash
set -e

PKG=jlcarveth-dev
SVCUSER=jlcarveth-dev
CFGDIR=/etc/$PKG
SHAREDIR=/usr/share/$PKG
NGINXAVAIL=/etc/nginx/sites-available/$PKG
NGINXENABLED=/etc/nginx/sites-enabled/$PKG

case "$1" in
    configure)
        # Create system user if it doesn't exist
        if ! id "$SVCUSER" &>/dev/null; then
            adduser --system --group --no-create-home --shell /usr/sbin/nologin "$SVCUSER"
        fi

        # Set ownership on config directory
        chown root:"$SVCUSER" "$CFGDIR"
        chmod 750 "$CFGDIR"

        # Create environment file from template if it doesn't exist
        if [ ! -f "$CFGDIR/environment" ]; then
            install -m 640 -o root -g "$SVCUSER" "$SHAREDIR/jlcarveth-dev.environment.template" "$CFGDIR/environment"
        fi

        # Install nginx configuration using ucf for upgrade-safe merging
        ucf --three-way "$SHAREDIR/jlcarveth-dev.nginx" "$NGINXAVAIL"
        ucfr "$PKG" "$NGINXAVAIL"
        chmod 644 "$NGINXAVAIL"

        # Enable nginx site if not already enabled
        if [ ! -e "$NGINXENABLED" ]; then
            ln -s "$NGINXAVAIL" "$NGINXENABLED"
        fi

        # Test and reload nginx
        if nginx -t 2>/dev/null; then
            systemctl reload nginx || true
        else
            echo "WARNING: nginx configuration test failed. Please check $NGINXAVAIL"
        fi

        echo ""
        echo "======================================================"
        echo " $PKG installed successfully"
        echo "======================================================"
        echo ""
        echo "Configuration: $CFGDIR/environment"
        echo ""
        echo "To start the service:"
        echo "  systemctl enable --now $PKG"
        echo ""
        echo "To check status:"
        echo "  systemctl status $PKG"
        echo "======================================================"
        echo ""
        ;;
esac

# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	# The following line should be removed in trixie or trixie+1
	deb-systemd-helper unmask 'jlcarveth-dev.service' >/dev/null || true

	# was-enabled defaults to true, so new installations run enable.
	if deb-systemd-helper --quiet was-enabled 'jlcarveth-dev.service'; then
		# Enables the unit on first installation, creates new
		# symlinks on upgrades if the unit file has changed.
		deb-systemd-helper enable 'jlcarveth-dev.service' >/dev/null || true
	else
		# Update the statefile to add new symlinks (if any), which need to be
		# cleaned up on purge. Also remove old symlinks.
		deb-systemd-helper update-state 'jlcarveth-dev.service' >/dev/null || true
	fi
fi
# End automatically added section
# Automatically added by dh_installsystemd/13.24.2
if [ "$1" = "configure" ] || [ "$1" = "abort-upgrade" ] || [ "$1" = "abort-deconfigure" ] || [ "$1" = "abort-remove" ] ; then
	if [ -d /run/systemd/system ]; then
		systemctl --system daemon-reload >/dev/null || true
		if [ -n "$2" ]; then
			_dh_action=restart
		else
			_dh_action=start
		fi
		deb-systemd-invoke $_dh_action 'jlcarveth-dev.service' >/dev/null || true
	fi
fi
# End automatically added section


exit 0
